# syborg-backup@.service
[Unit]
Description=run syborg backup

[Service]
Type=oneshot
ExecStart=/usr/bin/syborg backup %i
